# cron.d/Warning_Autotweet
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
HOME=/home/kouji/public_html/html/AutoTweet/
MAILTO=""

# m h dom mon dow user	command
#-----------------------------------------------------------------------------------------------------
30 10-20/2 * * *    	root `php bot-tenki-now2.php 三重県伊勢市 34.4550144 136.7236043 伊勢神宮`
30 11-21/2 * * *	root `php bot-tenki-now2.php 東京都杉並区 35.6828172 139.6359633 高千穂大学`
#-----------------------------------------------------------------------------------------------------
30 10-21 * * *		root `php EasyBotter/bot.php`
01 10-21 * * *		root `php bot-Denryoku-now.php NORMAL Linux`
31 10-21 * * *		root `php bot-Denryoku-now.php CHECK 高千穂大学`
#-----------------------------------------------------------------------------------------------------
